# flint.service
[Unit]
Description=Flint personal AI agent daemon
After=network.target signal-cli.service
Requires=signal-cli.service

[Service]
Type=simple
User=flint
EnvironmentFile=/etc/flint/flint.env
ExecStart=/usr/local/bin/flint --config /etc/flint/config.toml
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/flint
StandardOutput=journal
SyslogIdentifier=flint

[Install]
WantedBy=multi-user.target
